# 02-dev-hpc/03-ubuntu/base/supervisord.conf
#
# Supervisor 配置文件
#
# 功能:
#   用于管理容器内的后台服务，如此处的 SSH 服务。
#   确保即使主进程退出，服务也能自动重启。
#

[supervisord]
# 不以守护进程模式运行，因为 Docker 需要一个前台进程
nodaemon=true

[program:sshd]
# 定义要管理的 sshd 服务
# -D 选项使 sshd 在前台运行，而不是作为守护进程
command=/usr/sbin/sshd -D
# 服务启动时自动运行
autostart=true
# 服务退出后自动重启
autorestart=true
# 将服务的标准错误输出重定向到标准输出
redirect_stderr=true
# 将标准输出重定向到 stdout，以便 `docker logs` 可以捕获
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
